# systemd/vlog-worker.service
[Unit]
Description=VLog Transcoding Worker
After=network.target mnt-nas.mount
Wants=mnt-nas.mount

[Service]
Type=simple
User=damen
Group=damen
WorkingDirectory=/home/damen/vlog
ExecStart=/home/damen/vlog/venv/bin/python worker/transcoder.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
