# systemd/vlog-transcription.service
[Unit]
Description=VLog Transcription Worker
After=network.target mnt-nas.mount
Wants=mnt-nas.mount

[Service]
Type=simple
User=damen
Group=damen
WorkingDirectory=/home/damen/vlog
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/damen/vlog/venv/bin/python worker/transcription.py
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target
