# systemd/vlog-admin.service
[Unit]
Description=VLog Admin API
After=network.target mnt-nas.mount
Wants=mnt-nas.mount

[Service]
Type=simple
User=damen
Group=damen
WorkingDirectory=/home/damen/vlog
ExecStart=/home/damen/vlog/venv/bin/python -m uvicorn api.admin:app --host 0.0.0.0 --port 9001
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
