# systemd/vlog-worker.service
[Unit]
Description=VLog Transcoding Worker
After=network.target mnt-nas.mount
Wants=mnt-nas.mount

[Service]
Type=simple
User=damen
Group=damen
WorkingDirectory=/home/damen/vlog
Environment=PYTHONPATH=/home/damen/vlog/venv/lib/python3.9/site-packages
ExecStart=/usr/bin/python3 worker/transcoder.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
